# service/tor/run
#!/bin/sh

sh /etc/service/setup --passwd s3cret
tor -f /etc/service/tor/torrc