# service/privoxy/run
#!/bin/sh

privoxy --no-daemon /etc/service/privoxy/config
