# service/privoxy/run
#!/bin/sh
sh /etc/service/setup --address 0.0.0.0
privoxy --no-daemon /etc/service/privoxy/config
